# systemd/weightwise.service
[Unit]
Description=WeightWise Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/weightwise
EnvironmentFile=/root/weightwise/.env
ExecStart=/root/weightwise/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
